# cookbooks/node-lessons/templates/default/lessonsforlife.conf.erb
#!upstart
# lessonsforlife.conf
# Generated by Chef for <%= node['hostname'] %>.
# Local modifications will be overwritten.
# This script is to be placed in /etc/init to work with upstart.
#
# Internally the 'initctl' command is used to manage:
# initctl help
# initctl status lessonsforlife
# initctl reload lessonsforlife
# initctl start lessonsforlife

# modified from http://clock.co.uk/tech-blogs/upstart-and-nodejs

description "node.js server"
author      "Vincent Mac <vincent@avant.io>"
version     "1.0.1"


start on (local-filesystems and net-device-up)
stop on shutdown

instance "Node - lessonsforlife"

# Automatically Respawn
respawn
respawn limit 5 60

# Increase max open files (ulimit) from default @ 1024
limit nofile 32768 32768

# env HOME=/root

script
  #set -e

  export HOME="/root"

  echo $$ > /var/run/lessonsforlife.pid

  # exec sudo -u www-data /usr/bin/node /var/www/lessonsforlifeproject.com/app.js 2>&1 >> /var/log/lessonsforlifeproject.com
  # exec sudo -u <%= @user %> NODE_ENV=production <%= @node_bin %> <%= @app %> 2>&1 >> <%= @log %>
  exec sudo -u <%= @user %> NODE_ENV=production <%= @node_bin %> <%= @app %> >> <%= @log %> 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> <%= @log %>
end script

pre-stop script
    rm /var/run/lessonsforlife.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> <%= @log %>
end script